# forum/gunicorn.conf
description "Gunicorn application server handling forum"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
setuid ubuntu
setgid www-data
chdir /home/ubuntu/BESTHack-Dining-Room/forum
exec /home/ubuntu/venv/bin/gunicorn --workers 3 --bind localhost:8000 forum.wsgi:application
